# .config/systemd/user/ipo.service
[Unit]
Description=Automatic IPO

[Service]
Type=simple
ExecStart=/home/ayush/Code/Python/MeroShare-IPO/ipo.py --default True
TimeoutSec=300

